# images/backup/files/mdadm.cron
0 9 1 * * root /usr/sbin/check_array.sh

